# dream-server/opencode/opencode-web.service
[Unit]
Description=OpenCode Web UI
Documentation=https://opencode.ai/docs
After=network.target

[Service]
Type=simple
WorkingDirectory=__HOME__
ExecStart=__HOME__/.opencode/bin/opencode serve --port 3003 --hostname 127.0.0.1
Restart=on-failure
RestartSec=5

# Environment
Environment=HOME=__HOME__
Environment=PATH=__HOME__/.opencode/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-__HOME__/dream-server/.env

# Hardening
NoNewPrivileges=true
ReadWritePaths=__HOME__
PrivateTmp=true

[Install]
WantedBy=default.target
