# archive/guardian/guardian.service
[Unit]
Description=Guardian — Self-Healing Process Watchdog
Documentation=file:///etc/guardian/guardian.conf
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/guardian.sh
Restart=always
RestartSec=10

# Guardian runs as root; unprivileged users cannot kill or modify it
ProtectSystem=strict
# IMPORTANT: ProtectSystem=strict blocks ALL writes outside these paths.
# You MUST add every directory Guardian needs to read, backup, or restore.
# Without this, Guardian will silently fail to protect files in unlisted dirs.
# Example: ReadWritePaths=/var/log /var/lib/guardian /tmp /home/deploy /opt/myapp
ReadWritePaths=/var/log /var/lib/guardian /tmp
PrivateTmp=true

# Environment
Environment=GUARDIAN_CONF=/etc/guardian/guardian.conf

[Install]
WantedBy=multi-user.target
