# archive/systemd/llm-cold-storage.service
[Unit]
Description=LLM Cold Storage — Archive idle HuggingFace models
Documentation=https://github.com/Light-Heart-Labs/Lighthouse-AI

[Service]
Type=oneshot
ExecStart=%h/Lighthouse-AI/scripts/llm-cold-storage.sh --execute
Environment=HF_CACHE=%h/.cache/huggingface/hub
Environment=COLD_DIR=%h/llm-cold-storage
Environment=LOG_FILE=%h/.local/log/llm-cold-storage.log

# Run as the user who owns the HF cache
User=%i
