# archive/systemd/llm-cold-storage.timer
[Unit]
Description=Run LLM Cold Storage daily at 2am
Documentation=https://github.com/Light-Heart-Labs/Lighthouse-AI

[Timer]
OnCalendar=*-*-* 02:00:00
RandomizedDelaySec=900
Persistent=true

[Install]
WantedBy=timers.target
