# archive/systemd/openclaw-session-cleanup.timer
[Unit]
Description=Lighthouse AI - Session Cleanup every __INTERVAL__ minutes

[Timer]
OnBootSec=__BOOT_DELAY__min
OnUnitActiveSec=__INTERVAL__min
Persistent=true

[Install]
WantedBy=timers.target
