# archive/systemd/token-spy@.service
[Unit]
Description=Lighthouse AI - Token Spy (%i)
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__INSTALL_DIR__
EnvironmentFile=__INSTALL_DIR__/.env
Environment=AGENT_NAME=%i
ExecStart=/usr/bin/python3 -m uvicorn main:app --host __HOST__ --port __PORT__ --log-level warning
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
