# archive/systemd/vllm-tool-proxy.service
[Unit]
Description=Lighthouse AI - vLLM Tool Call Proxy
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=__USER__
ExecStart=/usr/bin/python3 __INSTALL_DIR__/vllm-tool-proxy.py --port __PROXY_PORT__ --vllm-url __VLLM_URL__
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
