# archive/configs/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway (Local vLLM)
After=network-online.target
Wants=network-online.target

[Service]
# Update the path to your openclaw installation
ExecStart=/usr/bin/node /usr/lib/node_modules/openclaw/dist/index.js gateway --port 18791
Restart=always
RestartSec=5
KillMode=process

# IMPORTANT: Set HOME to the directory containing your .openclaw/ config folder.
# OpenClaw reads ~/.openclaw/openclaw.json relative to HOME.
Environment=HOME=%h
Environment=OPENCLAW_GATEWAY_PORT=18791
Environment=OPENCLAW_SERVICE_KIND=gateway

[Install]
WantedBy=default.target
