# dream-server/scripts/systemd/memory-shepherd-workspace.timer
[Unit]
Description=Memory Shepherd — Workspace Files Timer (60s)

[Timer]
OnBootSec=20s
OnUnitActiveSec=60s
AccuracySec=5s

[Install]
WantedBy=timers.target
