# dream-server/scripts/systemd/memory-shepherd-memory.service
[Unit]
Description=Memory Shepherd — MEMORY.md Baseline Reset

[Service]
Type=oneshot
ExecStart=%h/dream-server/memory-shepherd/memory-shepherd.sh dream-agent-memory
