# dream-server/scripts/systemd/memory-shepherd-memory.timer
[Unit]
Description=Memory Shepherd — MEMORY.md Timer (3h)

[Timer]
OnBootSec=5min
OnCalendar=*-*-* 00/3:00:00
RandomizedDelaySec=5min
Persistent=true

[Install]
WantedBy=timers.target
