# dream-server/scripts/systemd/memory-shepherd-workspace.service
[Unit]
Description=Memory Shepherd — Workspace Files (AGENTS.md + TOOLS.md)

[Service]
Type=oneshot
ExecStart=%h/dream-server/memory-shepherd/memory-shepherd.sh dream-agent-agents
ExecStart=%h/dream-server/memory-shepherd/memory-shepherd.sh dream-agent-tools
